# scripts/config-updater.service
[Unit]
Description=CMBAgent Config Updater
After=network.target
Wants=network.target

[Service]
Type=simple
User=cmbagent
Group=cmbagent
WorkingDirectory=/opt/cmbagent-info
ExecStart=/usr/bin/python3 /opt/cmbagent-info/scripts/update-config-from-contexts.py --continuous --interval 600
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/cmbagent-info
Environment=GITHUB_TOKEN=your_github_token_here

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cmbagent-info

[Install]
WantedBy=multi-user.target 